# dev-cleanup/.config/systemd/user/dev-cache-cleanup.service
[Unit]
Description=Dev tool cache cleanup (pip, cargo, npm)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dev-cache-cleanup.sh
Nice=19
